# config_files/supervisord-uploadwatcher.conf
[supervidord]
nodaemon=true

[program:uploadwatcher]
command=/uploadwatcher.sh
